# itamae/cookbooks/env-variable/templates/etc/profile.d/variable.sh.erb
export SECRET_KEY_BASE=<%= Global.secret.rails.secret_key_base %>
export IKYUSAN_DATABASE_HOST=<%= Global.secret.database.ikyusan_database_host %>
export IKYUSAN_DATABASE_USER=<%= Global.secret.database.ikyusan_database_user %>
export IKYUSAN_DATABASE_PASSWORD=<%= Global.secret.database.ikyusan_database_password %>
export IKYUSAN_SLACK_WEBHOCK_URL=<%= Global.secret.slack.webhook_url %>
export NEW_RELIC_LICENSE_KEY=<%= Global.secret.newrelic.license_key %>
export IKYUSAN_DEFAULT_ICON_URL=<%= Global.secret.profile.defalut_icon_url %>
